# build/conf/upstart-latest_pledgers.conf
# at /etc/init/latest_pledgers.conf

description "Gunicorn application server running the latest_pledgers server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /opt/dxe/latest_pledgers/run_latest_pledgers.sh
