# build/conf/upstart-facebook.conf
# at /etc/init/facebook.conf

description "Gunicorn application server running the facebook data server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /opt/dxe/facebook/run_facebook.sh
