# build/conf/upstart-health.conf
# at /etc/init/health.conf

description "Gunicorn application server running the health check server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /opt/dxe/health/run_health.sh
